# services/davos/run
#!/bin/bash
exec /sbin/setuser abc java -jar /app/davos.jar 
